# supervisord.conf
; /etc/supervisor/conf.d/supervisord.conf

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root
strip_ansi=false

; Use ENV_API_APP coming from Docker ENV (defaults to app.main:app)
[program:api]
environment=API_APP="%(ENV_API_APP)s"
command=/bin/sh -lc "/venv/bin/python -m uvicorn $API_APP --host 0.0.0.0 --port 8000 --workers 1"
directory=/srv/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:web]
; Next standalone entry produced by output:'standalone'
command=node server.js
directory=/srv/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
